# tvcontrold/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /fs/jb/Library/LaunchDaemons/com.nito.tvcontrolservice.plist
