# tvcontrold/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/com.nito.tvcontrolservice.plist
